# slcan-tunnel-server.socket
[Unit]
Description=CAN tunneled over network

[Socket]
ListenStream=43410
Accept=true
MaxConnections=9
KeepAlive=true
#NoDelay=true

