# systemd/slcan-tunnel-server.socket
[Unit]
Description=CAN tunneled over network

[Socket]
ListenStream=6160
Accept=true
MaxConnections=9
KeepAlive=true
#NoDelay=true

